# sv/redis/run
#!/bin/bash

#sv start socklog-unix || exit 1

source /etc/envvars
env

mkdir -p /data

exec 2>&1
exec redis-server /etc/redis.conf
